# logforwarder.ini
[program:logforwarder-syslog]
command=/home/abe/anaconda3/bin/python /home/abe/work/log-forwarder/forwarder.py -s -c /home/abe/work/log-forwarder/config.yml
user=root
autorestart=true
stdout_logfile=/var/log/supervisor/logforwarder-syslog.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
redirect_stderr=true

[program:logforwarder-xflow]
command=/home/abe/anaconda3/bin/python /home/abe/work/log-forwarder/forwarder.py -f -c /home/abe/work/log-forwarder/config.yml
user=root
autorestart=true
stdout_logfile=/var/log/supervisor/logforwarder-xflow.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
redirect_stderr=true

[program:logforwarder-trap]
command=/home/abe/anaconda3/bin/python /home/abe/work/log-forwarder/forwarder.py -t -c /home/abe/work/log-forwarder/config.yml
user=root
autorestart=true
stdout_logfile=/var/log/supervisor/logforwarder-trap.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_capture_maxbytes=1MB
redirect_stderr=true
